# dist/zap.zapd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- macOS LaunchAgent: one zapd per user login, always up, self-binds the
     shared socket ~/.zap/run/zapd.sock. Install:
       cp dist/zap.zapd.plist ~/Library/LaunchAgents/
       launchctl load ~/Library/LaunchAgents/zap.zapd.plist -->
<plist version="1.0">
<dict>
  <key>Label</key><string>zap.zapd</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/zapd</string>
    <string>--log</string>
    <string>info</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/tmp/zapd.log</string>
  <key>StandardErrorPath</key><string>/tmp/zapd.log</string>
</dict>
</plist>
